# rootfs/etc/systemd/system/kresd@.service.d/override.conf
[Service]
LimitNOFILE=65535
EnvironmentFile=/etc/default/antizapret