# rootfs/etc/systemd/system/kresd@.service.d/override.conf
[Service]
LimitNOFILE=65535
PassEnvironment=DNS DNS_RU