# rootfs/etc/systemd/system/antizapret-update.timer
[Unit]
Description=Run antizapret-vpn update every 6h

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target